# src/deploy/gcp-vm/collector.service
[Unit]
Description=Turkey Business Activity — footfall collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__/src
Environment=PYTHONUNBUFFERED=1
# Two OpenMP threads matches the e2-micro's 2 shared vCPUs; torch's default
# oversubscribes and the context-switch thrash slows every YOLO pass down.
Environment=OMP_NUM_THREADS=2
# glibc grows one malloc arena per thread by default; on a threaded python
# process that quietly costs 50-150 MB of RSS - real money on a 1 GB host
# that was observed being taken down by the KERNEL oom-killer at 696M peak.
Environment=MALLOC_ARENA_MAX=2
Environment=FIREBASE_CREDENTIALS=__SA_PATH__
Environment=FIREBASE_STORAGE_BUCKET=__STORAGE_BUCKET__
# OSNet re-identification embedder. When the ONNX file is present the
# collector + dashboard-server use semantic identity matching (~5-10 ms
# per crop on CPU); when it's absent, both fall back to the HSV histogram
# transparently. Run `bash tools/setup_reid.sh` after `git pull` to fetch
# the model - it's a one-shot ~5 MB download and safe to skip if you're
# happy with the color-signature baseline.
Environment=REID_MODEL=__INSTALL_DIR__/src/data/osnet_x0_25_msmt17.onnx
ExecStart=__INSTALL_DIR__/src/.venv/bin/python -m app.collector --interval 40 --imgsz 512 --burst 2 --burst-stride 13
# --interval 40s: 5 slots x 3 writes/round - the effective cadence is bounded
# by the ROUND time on this host, not by the flag.
#
# --imgsz 512 --burst 2 --burst-stride 13 (2026-07, the 5-camera + OSNet era):
#   * measured on the live VM: 640/burst-3 with 5 cams + the OSNet embedder
#     pushed RSS past MemoryMax and the service entered a 4-minute oom-kill
#     loop (boot -> RSS climb -> SIGKILL -> RestartSec, repeat) while rounds
#     stretched past 7 minutes;
#   * 512 cuts YOLO FLOPs ~36% vs 640 and shrinks peak activations;
#   * burst 2 @ stride 13 keeps two frames ~0.5s apart - the median still
#     kills single-frame flicker and the speed estimator still gets its
#     two points - while dropping a third of the model passes and ~2/3 of
#     the H.264 decode work (detect_core also now picks the ~720p rendition
#     of the master playlist instead of decoding 1080p for a 512px input).
# On a host with more RAM (e2-small = 2 GB and up), restore
# `--imgsz 640 --burst 3 --burst-stride 25` for the last bit of recall.
Restart=always
RestartSec=15
# Memory sizing for the e2-micro (1 GB, Always Free) this project runs on.
# High/Max raised from 700/850 after the oom-kill loop above: the OSNet
# session + fifth camera legitimately need the headroom, and the flag
# changes on ExecStart pull peak RSS back down at the same time. If the
# journal shows oom-kills again, drop back to `--weights yolov8n.pt` or
# move to e2-small.
MemoryHigh=760M
MemoryMax=900M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
